# config/airootfs/etc/systemd/system/ntpd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ntpd -qg -u ntp:ntp
